# automation/ansible/roles/socat/templates/socat_instances/socat-7444d465-dc0a-440e-9f18-e161178ab27c.service
[Unit]
Description=Socat Daemon Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP4-LISTEN:3434,fork,reuseaddr TCP4:159.223.129.223:3434
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
