# automation/ansible/roles/socat/templates/socat_instances/socat-5352deff-bba2-4a31-af8a-933623f7f78d.service
[Unit]
Description=Socat Daemon Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP4-LISTEN:all,fork,reuseaddr TCP4:reimage:all
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
