# automation/ansible/roles/socat/templates/socat_instances/socat-b8441519-49a2-4dc2-98f0-2037b92e9514.service
[Unit]
Description=Socat Daemon Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP4-LISTEN:ansible,fork,reuseaddr TCP4:output:ansible
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
